# files/nginx.ini
[program:nginx]
command=/usr/bin/pidproxy /var/run/nginx.pid /usr/sbin/nginx
autostart = true
